# LOOBins/launchctl.yml
name: launchctl
author: Josh Carullo
short_description: Interact with LaunchAgents and LaunchDaemons. 
full_description: launchctl can be used to load, start, stop, and unload macOS services. It is a command-line frontend to launchd.
created: 2023-05-27
example_use_cases:
- name: Use launchctl to execute an application 
  description: A oneliner that will load a plist as a LaunchAgent or LaunchDaemon, achieving persistence on a target machine. This command requires root privileges.
  code: sudo launchctl load /Library/LaunchAgent/com.apple.installer
  tactics:
  - Execution
  - Persistence
  tags:
  - bash
  - zsh
  - oneliner
- name: Persistent launch agent 
  description: Creation of a persistent launch agent called with $HOME/Library/LaunchAgents/com.apple.updates.plist
  code: launchctl load -w ~/Library/LaunchAgents/com.apple.updates.plist
  tactics:
  - Persistence
paths:
- /bin/launchctl
detections:
- name: LaunchAgents and LaunchDaemons must have a plist file on disk in the root, system, or user Library directory. Monitoring for plist's with executables located in /tmp or /Shared could identify suspicious applications. 
  url: N/A
- name: "Jamf Protect: Detect launchctl activity that unloads or bootsout specific service"
  url: https://github.com/jamf/jamfprotect/blob/main/custom_analytic_detections/launchctl_unload_and_bootout_events
resources:
- name: '20 Common Tools & Techniques used by macOS threat Actors & Malware' 
  url: https://www.sentinelone.com/labs/20-common-tools-techniques-used-by-macos-threat-actors-malware/
- name: 'Mitre Attack Technique: launchctl T1569'
  url: https://attack.mitre.org/techniques/T1569/001/
- name: 'MITRE ATT&CK T1543.001  Create or Modify System Process: Launch Agent '
  url: https://attack.mitre.org/techniques/T1543/001/
- name: 'Komplex OS X Trojan (Sofacy)'
  url: https://unit42.paloaltonetworks.com/unit42-sofacys-komplex-os-x-trojan/
